# deploy/kairos.service
# kairos.service — service systemd « utilisateur » lançant « Kairos »
# (FastAPI/Uvicorn) au démarrage de la machine.
#
# Installation (depuis la racine du dépôt) : make service
# À la main (depuis la racine du dépôt) :
#   mkdir -p ~/.config/systemd/user
#   sed "s#__PROJECT_DIR__#$(pwd)#g" deploy/kairos.service \
#     > ~/.config/systemd/user/kairos.service
#   systemctl --user daemon-reload
#   systemctl --user enable --now kairos.service
#   loginctl enable-linger "$USER"   # optionnel : démarrage au boot sans session
#
# __PROJECT_DIR__ doit être le chemin ABSOLU de la racine du dépôt : WorkingDirectory
# doit pointer dessus pour que ExecStart (chemin relatif du venv) fonctionne, et pour
# qu'un éventuel .env d'une ancienne installation soit repéré au premier démarrage
# (migration automatique et unique, voir app/settings_store.py). Les réglages et la
# base de tâches vivent ensuite dans le dossier de données de l'utilisateur système
# (voir la page /kairos/settings de l'application), pas dans __PROJECT_DIR__.
#
# Proxy sortant (réseau d'entreprise) : réglable directement depuis la page Réglages
# de l'application (section « Réseau ») — plus besoin de fichier .env.proxy séparé.
#
# Port 8001 (l'outil de pilotage occupe 8000 : les deux coexistent sur un poste).

[Unit]
Description=Kairos (FastAPI/Uvicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__PROJECT_DIR__
ExecStart=__PROJECT_DIR__/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8001
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
